# templates/kubernetes/certs.timer.tmpl
[Unit]
Description=Periodic Certificate Renewal for {{.Component}}

[Timer]
OnActiveSec=0s
OnUnitActiveSec=12h

[Install]
WantedBy=default.target
