# templates/kubernetes/kube-logrotate.service.tmpl
[Unit]
Description=Kubernetes log rotation

[Service]
Type=oneshot
ExecStart=-/usr/sbin/logrotate /etc/logrotate.conf

[Install]
WantedBy=multi-user.target
