# templates/ip4tables.service.tmpl
[Unit]
Description=Packet Filtering Framework IPv4
DefaultDependencies=no
After=systemd-sysctl.service
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/iptables-restore /home/core/ip4tables.rules
ExecReload=/usr/sbin/iptables-restore /home/core/ip4tables.rules
ExecStop=/usr/sbin/iptables --flush
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
