# templates/update-cluster.service.tmpl
[Unit]
Description=Update PRIVATECLUSTER iptables chain

[Service]
Type=oneshot
ExecStart=/home/core/bin/yard cluster update --discovery-url={{.DiscoveryUrl}}
